# skel/etc/network/if-up.d/000route
#!/bin/sh
route add default gw 192.17.239.1
